# scripts/fusion360-mcp.service.template
[Unit]
Description=Fusion 360 MCP Server (SSE)
After=network.target

[Service]
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{WORKING_DIR}}
EnvironmentFile=/etc/default/fusion360-mcp
ExecStart={{VENV_PYTHON}} server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
